# scripts/crontab
# modloop daemon
0,15,30,45 * * * * /diva1/home/modloop/scripts/Modloop_sub.pl
